# profile.d/flutter.sh
#!/bin/bash
add_to_path ~/lib/flutter/bin
