# profile.d/go.sh
#!/bin/bash
export GOPATH="$HOME/lib/go"
add_to_path "$GOPATH/bin"