# profile.d/go.sh
#!/bin/bash
set -o allexport
GOPATH="$HOME/lib/go"
PATH="$GOPATH/bin:$PATH"