# profile.d/ls_aliases.sh
#!/bin/bash
alias ls='ls --color'
alias ll='ls -lh --color'
alias lll='ls -lah --color'
