# profile.d/misc_aliases.sh
#!/bin/bash
alias c=clear
alias exut=exit
alias exity=exit
alias exiy=exit
is_installed thefuck && eval "$(thefuck --alias)"
