# profile.d/osx_android.sh
#!/bin/bash
[[ -d /usr/local/opt/android-sdk ]] && export ANDROID_HOME=/usr/local/opt/android-sdk