# profile.d/osx_brew.sh
#!/bin/bash

alias brew_cu='brew cu -y --cleanup -q'
