# profile.d/osx_bashcompletion.sh
#!/bin/bash
for f in /usr/local/etc/bash_completion.d/*; do . "$f"; done