# profile.d/py3_venv.sh
#!/bin/sh
set -o allexport
VIRTUAL_ENV="~/py3.5"
PATH="$VIRTUAL_ENV/bin:$PATH"
unset PYTHON_HOME
